# zealot-web-api.service
[Unit]
Description=Zealot Web API - REST endpoints for status/guestbook
After=network.target

[Service]
Type=simple
User=aday
Group=aday
ExecStart=/usr/bin/python3 /home/aday/.local/bin/zealot_web_api.py
Restart=on-failure
RestartSec=10
Environment=TERM=linux
WorkingDirectory=/home/aday

[Install]
WantedBy=multi-user.target
