# zealot-hangs.service
[Unit]
Description=Zealot Hangs Bot - #ZealHangs multi-bot social channel
After=network.target ngircd.service
Wants=ngircd.service

[Service]
Type=simple
User=aday
Group=aday
ExecStart=/usr/bin/python3 /home/aday/.local/bin/zealot_hangs.py
Restart=on-failure
RestartSec=10
Environment=TERM=linux
WorkingDirectory=/home/aday

[Install]
WantedBy=multi-user.target
